# crontab
* * * * * cd /app && /opt/bitnami/php/bin/php artisan schedule:run >> /dev/null 2>&1
